# deb-package/chroot-terminal_1.0/DEBIAN/postinst
#!/bin/bash
set -e

echo "=== Chroot Terminal System v1.1 kurulumu ==="

# ── Kernel Optimizasyonu ─────────────────────
echo "[1/5] Kernel parametreleri uygulanıyor..."
sysctl -p /etc/sysctl.d/99-pty-limits.conf 2>/dev/null || true

# ── Ulimit Ayarlari ──────────────────────────
echo "[2/5] Ulimit ayarları uygulandı (limits.d)"
# PAM limits modülü zaten aktif, dosya koyulması yeterli

# ── Grup ve Kullanici ────────────────────────
echo "[3/5] Sistem grubu oluşturuluyor..."
groupadd -f ogrenciler

# ── Systemd Servisleri ───────────────────────
echo "[4/5] Servisler aktifleştiriliyor..."
systemctl daemon-reload
systemctl enable --now chroot-pty-fix.service
systemctl enable --now chroot-cleanup.timer

# ── SSH Reload ───────────────────────────────
echo "[5/5] SSH servisi yeniden yükleniyor..."
systemctl reload sshd 2>/dev/null || systemctl reload ssh 2>/dev/null || true

echo ""
echo "=== Kurulum tamamlandı ==="
echo ""
echo "  Kernel Optimizasyonu:"
echo "    PTY max:        16384"
echo "    File descriptors: 262144"
echo "    TCP backlog:    4096"
echo "    Swappiness:     10"
echo ""
echo "  Ulimit (genel/root):"
echo "    nofile:         65536 / 131072"
echo "    nproc:          32768 / 65536"
echo ""
echo "  Ulimit (ogrenci):"
echo "    nproc:          256 (fork bomb koruması)"
echo "    nofile:         4096"
echo ""
echo "  SSH:"
echo "    MaxSessions:    80"
echo "    MaxStartups:    50:60:100"
echo "    Compression:    off (LAN performans)"
echo "    UseDNS:         off (hızlı bağlantı)"
echo ""
echo "  Servisler:"
echo "    chroot-pty-fix: boot'ta PTY onarımı"
echo "    chroot-cleanup: her 10dk temizlik"
echo ""
echo "  Kullanım: chroot-yonetici init|create|list|mount|cleanup|health"
