# deb-package/chroot-terminal_1.1/etc/systemd/system/chroot-pty-fix.service
[Unit]
Description=Chroot PTY Repair on Boot
After=network.target ssh.service
Before=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/chroot-yonetici repair
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
