# deb-package/chroot-terminal_1.10/etc/systemd/system/chroot-network-fallback.service
[Unit]
Description=Chroot Terminal Network Fallback (raw ip)
After=network-pre.target
Before=network.target network-online.target ssh.service sshd.service
Wants=network-pre.target

# Sadece NM başarısız olursa çalış
ConditionPathExists=!/run/NetworkManager/NetworkManager.pid

[Service]
Type=oneshot
RemainAfterExit=yes
# Statik IP ile arayüzü kaldır (dhclient gerektirmez)
ExecStart=/bin/bash -c '\
  for iface in $(ls /sys/class/net/ | grep -E "^(en|eth)"); do \
    ip link set "$iface" up 2>/dev/null; \
    ip addr add 10.211.55.28/24 dev "$iface" 2>/dev/null || true; \
    ip route add default via 10.211.55.1 dev "$iface" 2>/dev/null || true; \
  done'
SuccessExitStatus=0 1

[Install]
WantedBy=multi-user.target
