# deb-package/chroot-terminal_1.10/etc/systemd/system/chroot-network-fallback.service
[Unit]
Description=Chroot Terminal Network Fallback (raw ip + dhclient)
After=network-pre.target systemd-udevd.service
Before=network.target network-online.target ssh.service sshd.service
Wants=network-pre.target

# Sadece networkd VE NM ikisi de başarısız olursa çalış
ConditionPathExists=!/run/NetworkManager/NetworkManager.pid

[Service]
Type=oneshot
RemainAfterExit=yes
# Tüm ethernet arayüzlerini bul ve DHCP ile kaldır
ExecStart=/bin/bash -c '\
  for iface in $(ls /sys/class/net/ | grep -E "^(en|eth)"); do \
    ip link set "$iface" up 2>/dev/null; \
    dhclient -1 -timeout 15 "$iface" 2>/dev/null || true; \
  done'
# Başarısız olsa da boot'u engelleme
SuccessExitStatus=0 1

[Install]
WantedBy=multi-user.target
