# deb-package/chroot-terminal_1.9/DEBIAN/postinst
#!/bin/bash
set -e

echo "=== Chroot Terminal System v1.9 kurulumu ==="

# ── LXC/Konteyner Tespiti ve Optimizasyonu ─────────────────
_is_container() {
    # Sadece gerçek konteyner ortamlarını tespit et (VM'ler hariç)
    # LXC, Docker, systemd-nspawn gibi ortamları yakalar
    # Parallels, VMware, VirtualBox gibi VM'leri YAKALAMASIN

    # 1. /proc/1/environ'da container= varsa (LXC, systemd-nspawn)
    grep -qi 'container=lxc\|container=docker\|container=systemd-nspawn' /proc/1/environ 2>/dev/null && return 0

    # 2. Docker ortamı
    [ -f /.dockerenv ] && return 0

    # 3. systemd-detect-virt --container: SADECE "lxc", "docker", "systemd-nspawn" dönerse
    if command -v systemd-detect-virt >/dev/null 2>&1; then
        local ctype
        ctype=$(systemd-detect-virt --container 2>/dev/null || echo "none")
        case "$ctype" in
            lxc|lxc-libvirt|docker|podman|systemd-nspawn|openvz)
                return 0 ;;
        esac
    fi

    return 1
}

if _is_container; then
    echo "[1/7] LXC/konteyner tespit edildi — gereksiz servisler devre dışı bırakılıyor..."
    for svc in systemd-udevd.service multipathd.service open-vm-tools.service tmp.mount; do
        systemctl disable "$svc" 2>/dev/null || true
        systemctl mask "$svc" 2>/dev/null || true
    done
    for sock in systemd-udevd-kernel.socket systemd-udevd-control.socket; do
        systemctl disable "$sock" 2>/dev/null || true
        systemctl mask "$sock" 2>/dev/null || true
    done
    if systemctl is-active systemd-resolved >/dev/null 2>&1 || systemctl is-enabled systemd-resolved >/dev/null 2>&1; then
        systemctl disable --now systemd-resolved.service 2>/dev/null || true
        systemctl mask systemd-resolved.service 2>/dev/null || true
        echo "nameserver 8.8.8.8" > /etc/resolv.conf 2>/dev/null || true
    fi
    mkdir -p /etc/systemd/system.conf.d
    cat > /etc/systemd/system.conf.d/lxc-timeout.conf << 'TIMEOUTEOF'
[Manager]
DefaultTimeoutStartSec=10s
DefaultTimeoutStopSec=10s
DefaultDeviceTimeoutSec=10s
TIMEOUTEOF
else
    echo "[1/7] Bare-metal/VM tespit edildi — sistem servisleri korunuyor."
fi

# ── Kernel Optimizasyonu ─────────────────────
echo "[2/7] Kernel parametreleri uygulanıyor..."
sysctl -p /etc/sysctl.d/99-pty-limits.conf 2>/dev/null || true

# ── Ulimit Ayarlari ──────────────────────────
echo "[3/7] Ulimit ayarları uygulandı (limits.d)"

# ── Grup ve Kullanici ────────────────────────
echo "[4/7] Sistem grubu oluşturuluyor..."
groupadd -f ogrenciler

# ── Systemd Servisleri ───────────────────────
echo "[5/7] Servisler aktifleştiriliyor..."
systemctl daemon-reload
systemctl enable --now chroot-pty-fix.service
systemctl enable --now chroot-cleanup.timer

# ── SSH Reload ───────────────────────────────
echo "[6/7] SSH servisi yeniden yükleniyor..."
systemctl reload sshd 2>/dev/null || systemctl reload ssh 2>/dev/null || true

# ── NTP Senkronizasyonu ──────────────────────
echo "[7/7] NTP zaman senkronizasyonu ayarlanıyor..."
timedatectl set-ntp true 2>/dev/null || true
systemctl restart systemd-timesyncd 2>/dev/null || true
timedatectl set-timezone Europe/Istanbul 2>/dev/null || true

echo ""
echo "=== Kurulum tamamlandı ==="
echo ""
echo "  Kernel Optimizasyonu:"
echo "    PTY max:        16384"
echo "    File descriptors: 262144"
echo "    TCP backlog:    4096"
echo "    Swappiness:     10"
echo ""
echo "  Ulimit (genel/root):"
echo "    nofile:         65536 / 131072"
echo "    nproc:          32768 / 65536"
echo ""
echo "  Ulimit (ogrenci):"
echo "    nproc:          256 (fork bomb koruması)"
echo "    nofile:         4096"
echo ""
echo "  SSH:"
echo "    MaxSessions:    80"
echo "    MaxStartups:    50:60:100"
echo "    Compression:    off (LAN performans)"
echo "    UseDNS:         off (hızlı bağlantı)"
echo ""
echo "  NTP:"
echo "    Timezone:       Europe/Istanbul"
echo "    NTP Servers:    ntp.ubuntu.com, time.google.com, time.cloudflare.com"
echo ""
if _is_container; then
echo "  LXC Boot Optimizasyonu:"
echo "    udevd, multipathd, open-vm-tools, tmp.mount: devre dışı"
echo "    systemd-resolved: statik DNS (8.8.8.8)"
echo "    Timeout: 10s (start/stop/device)"
else
echo "  Ortam: Bare-metal/VM (sistem servisleri korundu)"
fi
echo ""
echo "  Servisler:"
echo "    chroot-pty-fix: boot'ta PTY onarımı"
echo "    chroot-cleanup: her 10dk temizlik"
echo ""
echo "  Kullanım: chroot-yonetici init|create|list|mount|cleanup|health"
